# packaging/postinstall.sh
#!/usr/bin/env bash
systemctl enable porkbun-ddns
systemctl start porkbun-ddns
