# packaging/preremove.sh
#!/usr/bin/env bash
systemctl stop porkbun-ddns
systemctl disable porkbun-ddns
